# kvantum-watcher.service
[Unit]
Description=Kvantum Theme to KDE Plasma Color Sync
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kvantum-theme-watcher.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
